# .config/autostart/loloswitcher.desktop
[Desktop Entry]
Exec=/bin/loloswitcher
Icon=
Name=loloswitcher
Path=
Terminal=False
Type=Application
